# conf/statsd.conf
[program:statsd]
command=/usr/bin/nodejs /opt/statsd/stats.js /opt/statsd/config.js 2>&1 >> /var/log/statsd.log
autostart=true
autorestart=true
startretries=15
