# scripts/graphite.conf
[program:graphite]
command=/usr/bin/python /opt/graphite/webapp/graphite/manage.py runserver --noreload 0.0.0.0:80
user=root
autostart=true
autorestart=true
startretries=15
